# ansible/roles/head/files/rc.local
#!/bin/sh -e

/usr/local/bin/nat.sh

exit 0
